# deploy/aimilivpn-api-relay.service
[Unit]
Description=AimiliVPN restricted VPNGate API relay
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/aimilivpn
ExecStart=/usr/bin/python3 /opt/aimilivpn/vpngate_api_relay.py
EnvironmentFile=/etc/default/aimilivpn-api-relay
Restart=on-failure
RestartSec=10
DynamicUser=yes
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes

[Install]
WantedBy=multi-user.target
